# openclaw.service
[Unit]
Description=OpenClaw Swing Trading Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/openclaw
ExecStart=/opt/openclaw/venv/bin/python main.py
Restart=on-failure
RestartSec=30
StandardOutput=append:/opt/openclaw/logs/service.log
StandardError=append:/opt/openclaw/logs/service.log
EnvironmentFile=/opt/openclaw/.env

[Install]
WantedBy=multi-user.target
